# services/bitcoin/run
#!/bin/bash

echo "Starting bitcoind"
exec /sbin/setuser abc bitcoind -conf=/config/bitcoin.conf -datadir=/config/.bitcoin -daemon
